# launchd/com.openanthropic.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openanthropic.web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ourhome/bin/openanthropic-web-service.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ourhome/workstation/opencode-home/openanthropic</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/openanthropic-web.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/openanthropic-web.err.log</string>
</dict>
</plist>
